# experiments/overlay/etc/s6-overlay/s6-rc.d/comfyui/run
#!/command/with-contenv sh

. /opt/venv/comfyui/bin/activate

exec s6-setuidgid user \
     s6-envdir ./env \
     python /workspace/apps/comfyui/main.py --listen localhost --port 42010 \
     --output-directory /workspace/outputs/comfyui \
     --multi-user \
     --disable-auto-launch
