# runpod-base/overlay/etc/s6-overlay/s6-rc.d/caddy/run
#!/command/execlineb -P

/usr/bin/caddy run --watch --environ --config /etc/caddy/Caddyfile
