# runpod-base/overlay/etc/s6-overlay/s6-rc.d/setup-workspace/up
bash -c /opt/scripts/setup-workspace.sh
